# megaton.cron
3 * * * * /home/ubuntu/megaton/services/englishpepper/static/external/rescale_images.sh
